# Assingment_3_p2/login_daemon.service
[Unit]
Description=Login Daemon

[Service]
ExecStart=/usr/sbin/login_daemon
Restart=always

[Install]
WantedBy=multi-user.target